# templates/supervisor/rclone-serve.conf
{% import (role_path + "/templates/rclone/rclone")|relpath(playbook_dir) as rclone with context -%}
{% for config_name, config in rclone_serve.config.items() %}
[program:rclone-serve-{{ config_name }}]
user={{ config.daemon_user|default(rclone_serve.defaults.daemon_user, true) }}
command={{ rclone.rclone(config, rclone_serve.defaults) }}
{% endfor %}
